# systemd/fand.service
[Unit]
Description=fand fan control daemon
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/fand
# SIGHUP = hot-reload /etc/fand/config.toml (systemctl reload fand)
ExecReload=/bin/kill -HUP $MAINPID
ExecStopPost=/usr/local/bin/fand --restore-auto
Restart=on-failure
RestartSec=2

# Group=fand makes /run/fand owned root:fand, so 0750 lets group members
# reach the socket (which the daemon itself sets to 0660).
Group=fand
RuntimeDirectory=fand
RuntimeDirectoryMode=0750
ProtectSystem=strict
ReadWritePaths=/etc/fand
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
# NOTE: do NOT set ProtectKernelTunables=yes — it blocks hwmon sysfs writes.
# NVML needs access to /dev/nvidia*; if DevicePolicy is tightened, allow those nodes.

[Install]
WantedBy=multi-user.target
